# mavlink-router.service.in
[Unit]
Description=MAVLink Router

[Service]
Type=simple
ExecStart=@BINDIR@/mavlink-routerd --syslog
Restart=on-failure

[Install]
WantedBy=multi-user.target
